# services/echoos-camera.service
[Unit]
Description=EchoOS Camera Auth
After=echoos-ai.service

[Service]
Type=simple
User=echoos
Group=echoos
WorkingDirectory=/opt/echoos
ExecStart=/usr/bin/python3 /opt/echoos/camera/camera_face_verify.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
